# kyte_apm.service
[Unit]
Description=Memory and CPU Usage Monitor

[Service]
ExecStart=/path/to/your/ec2-apm.sh
